# etc/supervisor/conf.d/rq_worker.conf
[program:rq_worker]
command=/Users/mrdiamond/miniconda3/envs/laue-portal-conda-env/bin/python -m laue_portal.processing.worker
directory=/Users/mrdiamond/Documents/laue-portal
autostart=true
autorestart=true
stdout_logfile=/Users/mrdiamond/Documents/laue-portal/etc/supervisor/logs/rq_worker_%(process_num)02d.log
stderr_logfile=/Users/mrdiamond/Documents/laue-portal/etc/supervisor/logs/rq_worker_%(process_num)02d_err.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
numprocs=2
process_name=%(program_name)s_%(process_num)02d
